# rclone-sync.service
[Unit]
Description=Sync files using rclone

[Service]
Type=oneshot
User=root
ExecStart=/root/backup/rclone_borg_sync.sh